# service-monitor.service.template
[Unit]
Description=Service Monitor Dashboard
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=__WORKDIR__/node_modules/.bin/tsx server/index.ts
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
# Allow tsx to resolve node (e.g. an nvm install)
Environment=PATH=__NODE_BIN_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
